# bin/php-fpm/supervisord.conf
[program:php-fpm]
command = /usr/local/sbin/php-fpm
autostart=true
autorestart=true

[program:schedule]
command =php artisan schedule:work
autostart=true
autorestart=true

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[include]
files = /etc/supervisor/conf.d/*.conf